# cookbooks/main/files/default/mongodump.cron
SHELL=/bin/bash

0 3 * * *   root    /usr/local/bin/mongobackup > /dev/null